# rootfs/etc/services.d/traefik-admin/run
#!/usr/bin/with-contenv bashio

/web/traefik-admin --ConfigPath /data/dynamic.d --WebRoot /web/webroot