# rootfs/etc/services.d/traefik/run
#!/usr/bin/with-contenv bashio

ACME_EMAIL=$(bashio::config 'email')
mkdir -p /data/dynamic.d

HTTPREQ_ENDPOINT=$(bashio::config 'httpreq.url')

LOGLEVEL=$(bashio::config 'loglevel')

INSECURE_SKIP_VERIFY=$(bashio::config 'insecureSkipVerify')

SKIP_VERIFY=""
if [[ "$INSECURE_SKIP_VERIFY" = "true" ]]; then
    SKIP_VERIFY="--serversTransport.insecureSkipVerify=true"
fi
if [ ! -z "$HTTPREQ_ENDPOINT" ]; then
    HTTPREQ_USERNAME=$(bashio::config 'httpreq.username')
    HTTPREQ_PASSWORD=$(bashio::config 'httpreq.password')
    traefik \
    --log.level=$LOGLEVEL \
    $SKIP_VERIFY \
    --providers.file.directory=/data/dynamic.d \
    --entryPoints.web.address=:80 \
    --entryPoints.websecure.address=:443 \
    --certificatesresolvers.http01.acme.email=$ACME_EMAIL \
    --certificatesresolvers.http01.acme.storage=/data/acme.json \
    --certificatesresolvers.http01.acme.httpchallenge.entrypoint=web \
    --certificatesresolvers.dns01.acme.email=$ACME_EMAIL \
    --certificatesresolvers.dns01.acme.dnschallenge.provider=httpreq \
    --certificatesresolvers.dns01.acme.dnschallenge.delaybeforecheck=1 \
    --certificatesresolvers.dns01.acme.storage=/data/acme.json
else
    traefik \
    --log.level=$LOGLEVEL \
    $SKIP_VERIFY \
    --providers.file.directory=/data/dynamic.d \
    --entryPoints.web.address=:80 \
    --entryPoints.websecure.address=:443 \
    --certificatesresolvers.http01.acme.email=$ACME_EMAIL \
    --certificatesresolvers.http01.acme.storage=/data/acme.json \
    --certificatesresolvers.http01.acme.httpchallenge.entrypoint=web
fi


